# deploy/egypt-research-collect.service
[Unit]
Description=Collect and index Egypt Research sources
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/erx
ExecStart=/usr/bin/docker compose --env-file deploy/production.env run --rm egypt-research /app/deploy/collect-and-index.sh
